# ffrn-lowmem-patches/files/etc/init.d/ffrn-lowmem
#!/bin/sh /etc/rc.common

START=11

start () {
    /usr/sbin/ffrn-lowmem
}

stop() {
    true
}

