# ffs-lowmem-patches/files/etc/init.d/ffs-lowmem
#!/bin/sh /etc/rc.common

START=11

start () {
    /usr/sbin/ffs-lowmem
}

stop() {
    true
}

