# ffs-watchdog/files/etc/init.d/ffs-watchdog
#!/bin/sh /etc/rc.common

START=97

start () {
<<<<<<< Updated upstream
    rm /tmp/.*.ffs
=======
    rm /tmp/.gw.ffs
>>>>>>> Stashed changes
    /etc/init.d/micrond enable
    /etc/init.d/micrond restart
}

stop() {
    true
}
