# ffs-watchdog/files/etc/init.d/ffs-watchdog
#!/bin/sh /etc/rc.common

START=97

start () {
    rm /tmp/.*.ffs
    /etc/init.d/micrond enable
    /etc/init.d/micrond restart
}

stop() {
    true
}
